# code-editor/rootfs/etc/s6-overlay/s6-rc.d/init-code-server/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Code Editor
# Sets up code-server.
# ==============================================================================
declare config_path
readonly CODE_EDITOR_DATA_PATH="/data/code-editor"
readonly CODE_EDITOR_EXTENSIONS_PATH="${CODE_EDITOR_DATA_PATH}/extensions"
readonly CODE_EDITOR_SETTINGS_PATH="${CODE_EDITOR_DATA_PATH}/User/settings.json"
readonly CODE_SERVER_CONFIG_PATH="${CODE_EDITOR_DATA_PATH}/code-server.yaml"

# Hashes of previous default settings files. If a user's settings still match
# one of these, it is safe to replace them with the current defaults.
readonly -a PREVIOUS_DEFAULT_CONFIG_HASHES=(
    9d149651f7e8e17a7a0cfa6e5c7670a969f1f2d838eaed75ea13fa4591225692b79378de53e86f162cf7a48434c3cc5701e4fc5270c4f0fadf0d35982da676cb
)

if bashio::config.has_value 'config_path'; then
    config_path=$(bashio::config 'config_path')
    if ! bashio::fs.directory_exists "${config_path}"; then
        bashio::exit.nok "Configured config path does not exist"
    fi
    # Prevent using root directory as workspace to avoid performance issues
    if [[ "${config_path}" == "/" ]]; then
        bashio::exit.nok \
            "Using the root directory (/) as workspace is not allowed. " \
            "This causes severe performance issues as VS Code will index " \
            "the entire filesystem. Please use /config or another specific " \
            "directory instead."
    fi
fi

# Ensure persistent data folder exists.
if ! bashio::fs.directory_exists "${CODE_EDITOR_DATA_PATH}"; then
    mkdir -p "${CODE_EDITOR_EXTENSIONS_PATH}" \
        || bashio::exit.nok "Could not create persistent storage folder."
fi

# Ensure user extensions folder exists
mkdir -p "${CODE_EDITOR_EXTENSIONS_PATH}"

# Clear stale extension scanner state from previous app versions.
rm -f "${CODE_EDITOR_EXTENSIONS_PATH}/.obsolete" \
    "${CODE_EDITOR_EXTENSIONS_PATH}/extensions.json"

# Install bundled extensions into the runtime extensions folder.
while read -r ext; do
    ext="${ext//$'\r'/}"
    extension="${ext%%#*}"
    version="${ext##*#}"
    vsix="/usr/local/share/code-server/vsix/${extension}-${version}.vsix"
    installed_path=""

    [[ -z "${extension}" ]] && continue

    if ! bashio::fs.file_exists "${vsix}"; then
        bashio::log.warning "Bundled extension package not found: ${vsix}"
        continue
    fi

    installed_path=$(find "${CODE_EDITOR_EXTENSIONS_PATH}" \
        -maxdepth 1 \
        -type d \
        -iname "${extension}-${version}" \
        -print \
        -quit)

    if [[ -n "${installed_path}" ]]; then
        bashio::log.info "Bundled extension already installed: ${extension} ${version}"
        continue
    fi

    bashio::log.info "Installing bundled extension: ${extension} ${version}"

    if ! code-server \
        --user-data-dir "${CODE_EDITOR_DATA_PATH}" \
        --extensions-dir "${CODE_EDITOR_EXTENSIONS_PATH}" \
        --config "${CODE_SERVER_CONFIG_PATH}" \
        --install-extension "${vsix}"; then
        bashio::log.warning "Failed installing bundled extension: ${extension} ${version}"
    fi
done < /root/vscode.extensions

# Sets up default user settings on first start.
if ! bashio::fs.file_exists "${CODE_EDITOR_SETTINGS_PATH}"; then
    mkdir -p "${CODE_EDITOR_DATA_PATH}/User" \
        || bashio::exit.nok "Could not create persistent storage folder."

    cp /root/.code-server/settings.json "${CODE_EDITOR_SETTINGS_PATH}"
fi

# Upgrade settings.json is still default from previous version.
current=$(sha512sum "${CODE_EDITOR_SETTINGS_PATH}"|cut -d " " -f 1)
if [[ " ${PREVIOUS_DEFAULT_CONFIG_HASHES[*]} " == *" ${current} "* ]]; then
    cp /root/.code-server/settings.json "${CODE_EDITOR_SETTINGS_PATH}"
fi
